# etc/profile.d/gpg-agent.sh
#!/bin/bash

gpgconf --launch gpg-agent
